# week3-linux-administration/day3-4-process-service-management/scripts/custom-systemd/backup-monitor.service
# File: `/etc/systemd/system/backup-monitor.service`

[Unit]
Description=Backup Monitoring Service
After=network.target

[Service]
Type=simple
User=backup-user
WorkingDirectory=/opt/backup-monitor
ExecStart=/usr/bin/python3 /opt/backup-monitor/monitor.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target